# deploy/systemd/chatting-worker.service
[Unit]
Description=Chatting Worker Service (BBMB Consumer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=edward
Group=edward
WorkingDirectory=/home/edward/chatting
EnvironmentFile=/home/edward/chatting/configs/chatting-worker.env
ExecStart=/usr/bin/python3 -m app.main_worker
Restart=always
RestartSec=5
TimeoutStopSec=20
SyslogIdentifier=chatting-worker

[Install]
WantedBy=multi-user.target
